# xvfb/rootfs/etc/init.d/xvfb
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Xvfb Framebuffer"
command="/usr/bin/Xvfb"
command_args="${command_args:-}"


depend() {
  provide x11
}
